# uwsgi-apps.conf
description "uWSGI server for Django projects"

start on runlevel [2345]
stop on runlevel [!2345]
respawn
exec uwsgi --emperor="/server/apps/*/uwsgi.ini"